# ops/launchd/com.ai4sci.monitor.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.ai4sci.monitor.daily</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>cd __PROJECT_ROOT__ &amp;&amp; __PYTHON_BIN__ scripts/run_monitor_cycle.py --mode daily-full</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>3</integer>
    <key>Minute</key>
    <integer>15</integer>
  </dict>

  <key>WorkingDirectory</key>
  <string>__PROJECT_ROOT__</string>

  <key>StandardOutPath</key>
  <string>/Users/__USER__/Library/Logs/AI4Sci/monitor-daily.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/__USER__/Library/Logs/AI4Sci/monitor-daily.err.log</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
